# src/hh_applicant_tool/telegram/autokus.service
[Unit]
Description=AutoKus HR Bot Daemon
After=network.target

[Service]
Type=simple
User=flyingkuskus
WorkingDirectory=/home/flyingkuskus/AutoKus
ExecStart=/home/flyingkuskus/AutoKus/.venv/bin/python /home/flyingkuskus/AutoKus/daemon.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
